# deploy/systemd/ingero.service
[Unit]
Description=Ingero GPU Observability Agent
After=network.target
Documentation=https://github.com/ingero/ingero

[Service]
Type=simple
ExecStart=/usr/local/bin/ingero agent --config /etc/ingero/ingero.yaml
Restart=on-failure
RestartSec=5
LimitMEMLOCK=infinity
# eBPF requires CAP_BPF + CAP_PERFMON (or root)
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
